# deploy/raspberry-pi/asv-stack.target
[Unit]
Description=ASV backend and dual-camera runtime stack
Wants=network-online.target
After=network-online.target
Requires=asv-dashboard.service asv-stream.service
